# 服务器运维/rc.local/rc.local
#!/bin/sh -e
# 在exit 0前面插入需要服务器启动时执行的命令

exit 0
